# deploy/hermes-bot.service
[Unit]
Description=Hermes Download Nexus - Telegram Bot
After=network.target
Wants=hermes-api.service

[Service]
Type=simple
User=hermes
Group=hermes
WorkingDirectory=/opt/hermes
EnvironmentFile=/opt/hermes/.env
ExecStart=/opt/hermes/target/release/hermes-bot
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hermes-bot

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/hermes
PrivateTmp=true

[Install]
WantedBy=multi-user.target
